# roles/nfs_server/templates/nfs_vlan.service.jinja2
Unit]
Description=nfs_vlan init for feeder
Wants=network-online.target
After=network.target network-online.target

[Service]
ExecStart=/opt/nfs/init_vlan.sh

[Install]
WantedBy=multi-user.target