# roles/nfs_mount/templates/nfs_vlan.service.jinja2
[Unit]
Description=init nfs_vlan for mediabox
Wants=network-online.target
After=network.target network-online.target

[Service]
ExecStart=/opt/nfs/init_vlan.sh

[Install]
WantedBy=multi-user.target
