# roles/remote/templates/rclone_cache.service.j2
# /etc/systemd/system/rclone_cache.service
#########################################################################
# Title:         Cloudbox: Rclone Cache Mount                           #
# Author(s):     JoseP                                                  #
# URL:           https://github.com/cloudbox/cloudbox                   #
# --                                                                    #
#         Part of the Cloudbox project: https://cloudbox.works          #
#########################################################################
#                   GNU General Public License v3.0                     #
#########################################################################

[Unit]
Description=Rclone Cache Mount
After=network-online.target

[Service]
User={{ user.name }}
Group={{ user.name }}
Type=notify
ExecStartPre=/bin/sleep 10
ExecStart={{ rclone_cache_mount_start_command }}
ExecStop={{ rclone_cache_mount_stop_command }}
Restart=on-abort
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

[Install]
WantedBy=default.target
